# roles/provision-ib-node-exporter/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter for InfiniBand
After=network.target

[Service]
User=your_non_root_user
ExecStart={{ go_bin_path }}/node_exporter-infiniband

[Install]
WantedBy=default.target

